# deploy/systemd-files/kumanodocs-backup.timer
[Unit]
Description=Kumano Dormitory Document System Backup Task Timer

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
